# etc/systemd/system/getty@tty1.service.d/docker-monitor.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -n -o "-f docker" --noclear %I 38400 linux
TTYVTDisallocate=no
