# etc/systemd/system/docker.service.d/default.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker -D -d -s aufs -H fd://
