# etc/systemd/system/getty@.service.d/docker-monitor.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noclear -a docker %I $TERM
